# deploy/systemd/substation-web-gateway.service
[Unit]
Description=Substation FastAPI ROS Web Gateway
After=network-online.target substation-core.service
Wants=network-online.target

[Service]
Type=simple
User=substation
Group=substation
WorkingDirectory=/opt/substation/current
Environment=ROS_LOCALHOST_ONLY=1
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/substation/current/scripts/deployment/substation-web-gateway --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/substation/sqlite /var/lib/substation/evidence /var/lib/substation/reports /var/lib/substation/diagnostics

[Install]
WantedBy=multi-user.target
